# systemd/wallgen.service
[Unit]
Description=Service for a new wallpaper

[Service]
Restart=on-failure
ExecStart=/opt/WallGen/bin/wallgen
